# backup/postgresql_2023-08-06835680/cron.d/clean_backups_cron
0 0 */7 * * postgres /etc/clean_backups.sh
